# dotfiles/.config/systemd/user/dotfiles_snapshot.service
[Unit]
Description=Create a dotfiles snapshot

[Service]
Type=oneshot
ExecStart=%h/dotfiles/bin/snapshot create

[Install]
WantedBy=default.target
