# dotfiles/.config/systemd/user/dotfiles_snapshot.timer
[Unit]
Description=Create dotfiles snapshot every 10 minutes

[Timer]
OnActiveSec=10min

[Install]
WantedBy=timers.target