# imwheel/.config/systemd/user/imwheel.service
[Service]
Environment="DISPLAY=:1"
ExecStart=/usr/bin/imwheel --kill --buttons '4 5' -d
Restart=always

[Install]
WantedBy=default.target
